# wordops-panel/etc/systemd/system/wordops-panel.service
[Unit]
Description=WordOps Panel - WordOps GUI
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/opt/wordops-panel
Environment="PATH=/opt/wordops-panel/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/wordops-panel/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
