# wordops-panel/DEBIAN/prerm
#!/bin/bash
set -e

echo "Stopping WordOps Panel service..."
systemctl stop wordops-panel.service
systemctl disable wordops-panel.service

# Remove the virtual environment to clean up space
rm -rf /opt/wordops-panel/venv

exit 0
